# deploy/systemd/chadow-checkers-ws.service
[Unit]
Description=Chadow.ru online checkers WebSocket hub
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@RUN_USER@
Group=@RUN_GROUP@
WorkingDirectory=@CHADOW_ROOT@/deploy/checkers-ws
Environment=CHECKERS_WS_PORT=8792
EnvironmentFile=-@CHADOW_ROOT@/deploy/checkers-ws/.env
ExecStart=@NODE_BIN@ @CHADOW_ROOT@/deploy/checkers-ws/server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=chadow-checkers-ws

[Install]
WantedBy=multi-user.target
